# cliente/getty-tty6.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/ssh -i CHAVE_PRIVADA_SEM_SENHA USUARIO@ENDERECO_DO_SERVIDOR -p PORTA
StandardInput=tty
StandardOutput=tty
